# sleep.service
[Unit]
Description=Script Daemon

[Service]
Type=simple
ExecStart=/tmp/sleep.sh
Restart=always

[Install]
WantedBy=default.target
